# not-needed/tools/eggs/live-usb/systemd/home-casparcg-exfat.mount.example
# Install as: /etc/systemd/system/home-casparcg-exfat.mount
# Replace UUID= with:  blkid /dev/sdXN
# Requires: exfat-fuse or kernel exfat driver (linux 5.4+).
#
# Mount options use uid/gid of the HighAsCG service user (default casparcg).

[Unit]
Description=HighAsCG exFAT data volume at /home/casparcg/exfat
DefaultDependencies=no
Conflicts=umount.target
Before=highascg.service highascg-exfat-sync.service
After=blk-availability.target
After=systemd-remount-fs.service

[Mount]
What=/dev/disk/by-uuid/00000000-0000-0000-0000-000000000000
Where=/home/casparcg/exfat
Type=exfat
Options=defaults,uid=1000,gid=1000,umask=002,nofail

[Install]
WantedBy=local-fs.target
