# tools/eggs/live-usb/nvidia-multi-driver/highascg-pick-nvidia.service
[Unit]
Description=HighAsCG: pick recommended NVIDIA driver on first boot
Documentation=file:///usr/local/sbin/highascg-pick-nvidia.sh
ConditionPathExists=!/var/lib/highascg/nvidia-installed
After=local-fs.target
Before=highascg.service display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/highascg-pick-nvidia.sh
RemainAfterExit=yes
TimeoutStartSec=900
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
